# systemd/jai-voice.service
[Unit]
Description=JAI Voice Assistant – always‑on LLM + mic + speaker
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/wiju/projects/jaihome
# Start the Docker compose stack (foreground)
ExecStart=/usr/bin/docker compose -f docker/compose.yml up
# Stop the stack on service stop
ExecStop=/usr/bin/docker compose -f docker/compose.yml down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
